# init-scripts/custom-start/custom.target
[Unit]
Description=Custom Target for skywire node
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
